# crontab
0 0 * * * root python3 /recycle.py >> /var/log/cron.log
